# t/dh_installinit/debian/foo.service
[Unit]
Description=A unit

[Service]
ExecStart=/bin/true
