# t/dh_installsystemd/debian/foo2.service
[Unit]
Description=Another unit

[Service]
ExecStart=/bin/true

[Install]
WantedBy=multi-user.target
