# t/dh_installsystemduser/debian/baz.user.service
[Unit]
Description=Baz User Unit

[Service]
ExecStart=/bin/true

[Install]
WantedBy=default.target
