# t/dh_installsystemduser/debian/foo.user.service
[Unit]
Description=Foo User Unit

[Service]
ExecStart=/bin/true

[Install]
WantedBy=default.target
